# thesis_record_app/config/deploy/systemd/thesis-record-web.service.example
[Unit]
Description=ThesisRecord Rails web process
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=thesis-record
WorkingDirectory=/var/www/thesis-record/thesis_record_app
EnvironmentFile=/etc/thesis-record/thesis-record.env
ExecStart=/bin/bash -lc 'bundle exec puma -C config/puma.rb'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
