# thesis_record_app/config/deploy/systemd/user/thesis-record-sidekiq-watchdog.service.example
[Unit]
Description=ThesisRecord Sidekiq active-state watchdog (user service)

[Service]
Type=oneshot
ExecStart=/bin/bash -lc 'systemctl --user is-active --quiet thesis-record-sidekiq.service'
