# thesis_record_app/config/deploy/systemd/user/thesis-record-web.service.example
[Unit]
Description=ThesisRecord Rails web process (user service)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/thesis-record/thesis_record_app
EnvironmentFile=%h/.config/thesis-record/thesis-record.env
ExecStart=/bin/bash -lc 'bundle exec puma -C config/puma.rb'
Restart=always
RestartSec=5
KillSignal=SIGTERM

[Install]
WantedBy=default.target
